# batteryMonitor/battery-check.service
[Unit]
Description=Controllo livello batteria

[Service]
Type=oneshot
ExecStart=%h/.local/bin/monitor_battery.sh
