# backupHDD/backup-hdd.service
[Unit]
Description=Pipeline Backup HDD USB

[Service]
Type=oneshot
ExecStart=/usr/bin/ionice -c 3 /usr/local/bin/backup_hdd.sh
User=noya
